# rpm/demo-libuv.service
[Unit]
Description=demo_libuv oneshot service

[Service]
Type=oneshot
ExecStart=/usr/bin/demo_libuv
